# deploy/pweb.service
[Unit]
Description=Gunicorn for PWeb (LyomasPharmacy)
After=network.target postgresql.service

[Service]
User=softwarehouse
Group=www-data
WorkingDirectory=/var/www/pweb
EnvironmentFile=/var/www/pweb/.env
ExecStart=/var/www/pweb/venv/bin/gunicorn config.wsgi:application \
    --workers 3 \
    --bind unix:/var/www/pweb/gunicorn.sock \
    --timeout 120 \
    --access-logfile /var/log/pweb/access.log \
    --error-logfile /var/log/pweb/error.log
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
